# etc/supervisord.d/diamond.ini
[program:diamond]
command=/root/bin/start_diamond.sh
redirect_stderr=true
stdout_syslog=true
