# etc/supervisord.d/setup.ini
[program:setup]
command=/usr/local/bin/qnib-setup.py
priority=1
startretries=0
startsecs=0
redirect_stderr=true
stdout_syslog=true

