# etc/supervisord.d/syslog-ng.ini
[program:syslog-ng]
command=/root/bin/start_syslogng.sh
redirect_stderr=true
stdout_syslog=true
#stdout_logfile = /var/log/supervisor/%(program_name)s.log
#stdout_capture_maxbytes=1MB
