# etc/supervisord.d/watchpsutil.ini
[program:watchpsutil]
command=/opt/qnib/bin/watch_psutil.py -l --loop
redirect_stderr=true
stdout_logfile=syslog
