# etc/supervisord.d/watchpsutil.ini
[program:watchpsutil]
command=/opt/qnib/bin/watch_psutil.py -l --loop
startsecs=1
exitcodes=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
